# systemd/sdwan-traffic-gen.service
[Unit]
Description=SD-WAN Traffic Generator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/sdwan-traffic-gen/traffic-generator.sh client01
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
